# tests/linux_repro/cron.d/malicious
* * * * * root /usr/bin/python3 /tmp/backdoor.py
